# launchagents/com.wezterm.mux-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wezterm.mux-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/wezterm-mux-server</string>
        <string>--config-file</string>
        <string>/Users/artharol/.wezterm.lua</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/Users/artharol/Library/Logs/wezterm-mux-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/artharol/Library/Logs/wezterm-mux-server.log</string>
</dict>
</plist>
